# resources/share/i3-xfce/roles/i3/files/autostart/xfce4-panel.desktop
[Desktop Entry]
Type=Application
Terminal=false
Name=XFCE Panel
Comment=Panel
Exec=xfce4-panel --disable-wm-check
OnlyShowIn=XFCE;
